# Builds/RobotSimulator/Linux x64/robotsimulator.service
[Unit]
Description=RobotSimulator

[Service]
ExecStart=/usr/sbin/RobotSimulator

[Install]
WantedBy=multi-user.target